# zpool-scrub@.timer
[Unit]
Description=Scrub ZFS pool semi-weekly

[Timer]
# Every other week ish
OnCalendar=*-*-1,15
Persistent=true

[Install]
WantedBy=timers.target
